# service/tead.service
[Unit]
Description=Tea Daemon, listens for remote clients informing us of tea status.

[Service]
Type=oneshot
ExecStart=/usr/bin/tea
